# contrib/systemd/check_vpn.service.d/limits.conf
[Service]
# Resource limits and runtime directory mode
RuntimeDirectoryMode=0750
LimitNOFILE=65536
LimitNPROC=512
# Consider adding LimitCORE=0 to disable core dumps if desired by packagers
